# crontab
0 0 * * 0 make --quiet --directory /opt/docker
